# app/views/cloud_model/guest/etc/systemd/system/nginx.service.d/certbot_init.conf.erb
[Service]
ExecStartPost=/usr/bin/certbot --nginx --agree-tos --email <%= CloudModel.config.admin_email %> --domains <%= guest.external_hostname %> -n