# app/views/cloud_model/guest/etc/systemd/system/nginx.service.d/fix_perms.conf.erb
[Service]
PermissionsStartOnly=true
ExecStartPre=/bin/chown -R www:www /var/www