# app/views/cloud_model/guest/etc/systemd/system/nginx.service.d/fix_perms.conf.erb
[Service]
PermissionsStartOnly=true
ExecStartPre=/bin/chown -R www:www /var/www
ExecStartPre=-/usr/bin/sudo -u www /var/www/rails/current/script/prestart