# app/views/cloud_model/guest/etc/systemd/system/certbot-renew.timer.erb
[Unit]
Description=Call certbot-renew on schedule

[Timer]
OnBootSec=180
# run 3 minutes after boot time
OnCalendar=03:00
AccuracySec=600

[Install]
WantedBy=timers.target