# app/views/cloud_model/guest/etc/systemd/system/certbot-renew.service.erb
[Unit]
Description=Renew certbot certificates

[Service]
ExecStart=/usr/bin/certbot --nginx renew