# app/views/cloud_model/guest/etc/systemd/system/cm_fix_permissions.service.erb
[Unit]
Description=Ensure some rights on starting the system
Before=network.service

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/share/cloud_model/fix_permissions.sh

[Install]
WantedBy=basic.target