# app/views/cloud_model/guest/etc/systemd/system/delayed_jobs.service.erb
[Unit]
<%- path = "/var/www/rails/current" %>
Description=Delayed Job Worker
After=syslog.target network.target
AssertPathExists=<%= path %>

[Service]
Type=forking
RemainAfterExit=no
WorkingDirectory=<%= path %>
PIDFile=/tmp/delayed_job.pid

Environment = RAILS_ENV=<%= model.passenger_env %>
ExecStart=<%= path %>/bin/delayed_job --pid_dir=/tmp/ start
ExecStop=<%= path %>/bin/delayed_job --pid_dir=/tmp/ stop
User=www
Group=www
UMask=0002

# Restart after 5 seconds delay
Restart=always
RestartSec=5

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=dj_worker

[Install]
WantedBy=multi-user.target