# app/views/cloud_model/guest/etc/systemd/system/dhclient.service.erb
[Unit]
Description=dhclient on all interfaces
Wants=network.target
Before=network.target

[Service]
ExecStart=/sbin/dhclient -4 -q

[Install]
WantedBy=multi-user.target