# app/views/cloud_model/guest/etc/systemd/system/mongodb.service.erb
[Unit]
Description=MongoDB schema-free document-oriented database
After=syslog.target network.target auditd.service

[Service]
Environment="LC_ALL=C"
ExecStart=/usr/bin/mongod --quiet --config /etc/mongodb.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
User=mongodb
Group=mongodb

[Install]
WantedBy=multi-user.target