# app/views/cloud_model/guest/etc/systemd/system/rake@.timer.erb
[Unit]
Description=Call rake on schedule

[Timer]
OnBootSec=15min    # run 15 minutes after boot time
OnCalendar=00:00
AccuracySec=10min

[Install]
WantedBy=timers.target