# app/views/cloud_model/guest/etc/systemd/system/rake@.timer.erb
[Unit]
Description=Call rake on schedule

[Timer]
OnBootSec=900
# run 15 minutes after boot time
OnCalendar=00:00
AccuracySec=600

[Install]
WantedBy=timers.target