# app/views/cloud_model/guest/etc/systemd/system/redis.service.erb
[Unit]
Description=A persistent key-value database
After=syslog.target network.target auditd.service

[Service]
ExecStart=/usr/sbin/redis-server /etc/redis.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
User=redis
Group=redis

[Install]
WantedBy=multi-user.target