# app/views/cloud_model/guest/etc/systemd/system/shinken-reactionner.service.erb
[Unit]
Description=Shinken Reactionner
After=syslog.target

[Service]
ExecStart=/usr/bin/shinken-reactionner -r -c /etc/shinken/daemons/reactionnerd.ini

[Install]
WantedBy=multi-user.target