# app/views/cloud_model/guest/etc/systemd/system/shinken-arbiter.service.erb
[Unit]
Description=Shinken Arbiter
After=syslog.target

[Service]
ExecStart=/usr/bin/shinken-arbiter -r -c /etc/shinken/shinken.cfg

[Install]
WantedBy=multi-user.target