# app/views/cloud_model/guest/etc/systemd/system/shinken-broker.service.erb
[Unit]
Description=Shinken Broker
After=syslog.target

[Service]
ExecStart=/usr/bin/shinken-broker -r -c /etc/shinken/daemons/brokerd.ini

[Install]
WantedBy=multi-user.target