# app/views/cloud_model/guest/etc/systemd/system/shinken-receiver.service.erb
[Unit]
Description=Shinken Receiver
After=syslog.target

[Service]
ExecStart=/usr/bin/shinken-receiver -r -c /etc/shinken/daemons/receiverd.ini

[Install]
WantedBy=multi-user.target