# app/views/cloud_model/guest/etc/systemd/system/shinken-scheduler.service.erb
[Unit]
Description=Shinken Scheduler
After=syslog.target

[Service]
ExecStart=/usr/bin/shinken-scheduler -r -c /etc/shinken/daemons/schedulerd.ini

[Install]
WantedBy=multi-user.target