# app/views/cloud_model/guest/etc/systemd/system/solr.service.erb
[Unit]
Description=SOLR Search Engine
After=network.target

[Service]
User=solr
Group=solr
ExecStart=/opt/solr/bin/solr start -f -m 2g -s /var/solr/solr -p <%= model.port %> -a "-Dsolr.data.dir=/var/solr/data"
ExecStop=/opt/solr/bin/solr stop -p 8080
Restart=always

[Install]
WantedBy=multi-user.target