# app/views/cloud_model/guest/etc/systemd/system/solr.service.erb
[Unit]
Description=SOLR Search Engine
After=network.target

[Service]
Environment=SOLR_JAVA_MEM="-Xms512m -Xmx<%= model.heap_size %>"
User=solr
Group=solr
PermissionsStartOnly=true
ExecStartPre=/bin/chown -R solr:solr /var/solr/data
ExecStart=/opt/solr/bin/solr start -f -s /var/solr/solr -p <%= model.port %> -a "-Dsolr.data.dir=/var/solr/data -Dsolr.data.logs=/var/solr/log -Xloggc:/var/solr/log/solr_gc.log"
ExecStop=/opt/solr/bin/solr stop -p 8080
Restart=always

[Install]
WantedBy=multi-user.target