# app/views/cloud_model/guest/etc/systemd/system/tomcat8.service.erb
[Unit]
Description=Tomcat 8 servlet container
After=network.target

[Service]
User=tomcat8
Group=tomcat8
Type=forking
PIDFile=/run/tomcat8/pid
ExecStart=/usr/sbin/tomcat8 start
ExecStop=/usr/sbin/tomcat8 stop
Restart=always

[Install]
WantedBy=multi-user.target