# app/views/cloud_model/guest/etc/systemd/system/tomcat-8.service.erb
[Unit]
Description=Tomcat 8 servlet container
After=network.target

[Service]
User=tomcat
Group=tomcat
ExecStart=/usr/sbin/tomcat-8 start
ExecStop=/usr/sbin/tomcat-8 stop
Restart=always

[Install]
WantedBy=multi-user.target