# app/views/cloud_model/guest/etc/systemd/system/tomcat-7.service.erb
[Unit]
Description=Tomcat 7 servlet container
After=network.target

[Service]
User=tomcat
Group=tomcat
ExecStart=/usr/sbin/tomcat-7 start
ExecStop=/usr/sbin/tomcat-7 stop
Restart=always

[Install]
WantedBy=multi-user.target