# app/views/cloud_model/host/etc/systemd/system/guest_zpool.service.erb
[Unit]
Description=Bootstrap zfs pool for guests
Before=network.target lxd.service
Conflicts=shutdown.target

[Service]
Type=oneshot
# Try to create guests pool; exists with error code if already exist
ExecStartPre=/sbin/zpool create guests mirror sda6 sdb6
# Init lxd if zpool create did not exit with error code aka create success
ExecStart=/usr/bin/lxd init --auto --storage-backend zfs --storage-pool guests

[Install]
WantedBy=basic.target

