# app/views/cloud_model/guest/etc/systemd/system/solr.service.d/fix_perms.conf.erb
[Service]
PermissionsStartOnly=true
ExecStartPre=/bin/chown -R solr:solr /var/solr