# app/views/cloud_model/guest/etc/systemd/system/neo4j.service.d/fix_perms.conf.erb
[Service]
PermissionsStartOnly=true
ExecStartPre=/bin/chown -R neo4j:adm /var/lib/neo4j