# app/views/cloud_model/host/etc/init.d/cloudmodel.erb
#!/sbin/runscript

depend() {
  before net
  use logger
}

start() {
  ebegin "Starting CloudModel"
  <%= firewall_worker.start_script %>
  eend $?
}

stop() {
  ebegin "Stopping CloudModel"
  <%= firewall_worker.stop_script %>
  eend $?
}
