# app/views/cloud_model/guest/etc/systemd/system/graphite-web.service.d/fix_perms.conf.erb
[Service]
PermissionsStartOnly=true
ExecStartPre=/bin/chown -R _graphite:_graphite /var/lib/graphite/