# app/views/cloud_model/guest/etc/systemd/system/forgejo.service.d/fix_perms.conf.erb
[Service]
PermissionsStartOnly=true
ExecStartPre=/bin/chown -R forgejo:forgejo /var/lib/forgejo